# DEBIAN/prerm
#!/bin/sh
# prerm script for sickragern6
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see /usr/share/doc/packaging-manual/

case "$1" in
	remove|upgrade|deconfigure)

		###
		# Disable our service
		###
		if [ `systemctl is-active fvapp-sickragern6` = "active" ]; then
			echo "stopping fvapp-sickragern6"
        		systemctl stop fvapp-sickragern6
		else 
        		echo "fvapp-sickragern6 not running"
		fi

		if [ -f /lib/systemd/system/fvapp-sickragern6.service ]; then
			echo "disable fvapp-sickragern6"
        		systemctl disable fvapp-sickragern6
		else
        		echo "fvapp-sickragern6 not enabled"
		fi

		###
		# Remove compiled python stuff
		###
		echo "remove compiled python stuff"
		find /apps/sickragern6 -name "*.pyc" -print0 | xargs -n1 -0 rm -f
		find /apps/sickragern6 -name "*.pyo" -print0 | xargs -n1 -0 rm -f
		echo "completed removing compiled python stuff"

		###
		# remove sickragern6 app
		###
		if [! -f /apps/sickragern6 ]; then
			echo "removing sickragern6 app folder"
                	rm -rf /apps/sickragern6
			event_push app readynasd '<delete-s resource-type="LocalApp" resource-id="LocalApp"><LocalApp appname="sickragern6" success="1"/></delete-s>' 0 0
			echo "sickragern6 Removed"
		else
			echo "Error sickragern6 directory still present"
		fi
	;;
	failed-upgrade)
	;;
	*)
		echo "prerm called with unknown argument \`$1'" >&2
        	exit 0
    	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.


exit 0

